# .config/systemd/user/weather-updater.service
[Unit]
Description=Update weather cache

[Service]
ExecStart=%h/.config/waybar/custom_modules/weather.sh

